# crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
* * * * * php -r "echo date(\DateTime::W3C), PHP_EOL;" > /proc/1/fd/1 2>/proc/1/fd/2
